# rc.local
#!/bin/sh -e
arp -f /etc/ethers
systemctl start ssh
exit 0
